# ecmwf/ecmwf.cron
15 * * * * /smartmet/bin/ingest-data.sh -m ecmwf
