# ukmo/ukmo.cron
15 * * * * /smartmet/bin/ingest-model.sh -m ukmo
