# wrf/wrf.cron
15 * * * * /smartmet/bin/ingest-model.sh -m wrf -a small
15 * * * * /smartmet/bin/ingest-model.sh -m wrf -a large
