# gsm/gsm.cron
15 * * * * /smartmet/bin/ingest-model.sh -m gsm
